# deploy/supervisor.conf
; Supervisord-Config fuer HostPoint FlexServer hpservices.
;
; Diesen Inhalt nach Anlegen des Service kopieren:
;   1. hpservices supervisord add auftragsverwaltung
;   2. nano ~/.services/supervisord/auftragsverwaltung/service.conf
;      → diesen Inhalt einfuegen, speichern
;   3. supervisorctl update
;   4. supervisorctl status auftragsverwaltung
;
; Logs anschauen:
;   tail -f ~/apps/auftragsverwaltung/logs/stdout.log
;   tail -f ~/apps/auftragsverwaltung/logs/stderr.log

[program:auftragsverwaltung]
command=/home/xahizivi/apps/auftragsverwaltung/deploy/run.sh
directory=/home/xahizivi/apps/auftragsverwaltung
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM

; Logging
stdout_logfile=/home/xahizivi/apps/auftragsverwaltung/logs/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/xahizivi/apps/auftragsverwaltung/logs/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
